# codex-terminal/rootfs/etc/services.d/codex/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

PORT=7681
SHELL_BIN="/bin/bash"

if ! command -v ttyd >/dev/null 2>&1; then
    bashio::log.fatal "ttyd is missing inside the container."
    exit 1
fi

if [ ! -x "${SHELL_BIN}" ]; then
    SHELL_BIN="/bin/sh"
fi

bashio::log.info "Starting Codex Terminal (ttyd) on port ${PORT}."
exec ttyd --writable -p "${PORT}" -i 0.0.0.0 "${SHELL_BIN}"
