# codex-terminal/rootfs/etc/services.d/codex/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

PORT=7681

if ! command -v ttyd >/dev/null 2>&1; then
    bashio::log.fatal "ttyd is missing inside the container."
    exit 1
fi

TTY_CMD="/bin/bash"
TTY_ARGS=""

if command -v codex >/dev/null 2>&1; then
    TTY_ARGS="-lc codex"
else
    bashio::log.warning "codex CLI not found in PATH; launching plain shell."
fi

bashio::log.info "Starting Codex Terminal (ttyd) on port ${PORT}."
if [ -n "${TTY_ARGS}" ]; then
    exec ttyd --writable -p "${PORT}" -i 0.0.0.0 "${TTY_CMD}" ${TTY_ARGS}
else
    exec ttyd --writable -p "${PORT}" -i 0.0.0.0 "${TTY_CMD}"
fi
